# heartbeat-service/debian/heartbeat.service
[Unit]
Description=Heartbeat Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/heartbeat-service
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=heartbeat

[Install]
WantedBy=multi-user.target